# supervisor/supervisord-gin.conf
[program:golang-dataapi]
command = /usr/local/go/src/DataApi.Go/main
stopsignal=QUIT
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true